# conf/novelsearch.tpl.conf
[program:novelsearch]
command={{PATH}}/novelsearch
autostart=true
autorestart=true
startsecs=10
#/etc/supervisor/novelsearch.conf