# scripts/novelsearch.conf
[program:novelsearch]
command=/Users/midoks/Desktop/gopath/src/github.com/midoks/novelsearch/novelsearch
autostart=true
autorestart=true
startsecs=10
#/etc/supervisor/novelsearch.conf
